# disk2iso-lib/service/disk2iso.service
[Unit]
Description=disk2iso - Automatische ISO Erstellung von optischen Medien
After=multi-user.target
Wants=systemd-udevd.service
After=systemd-udevd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/disk2iso.sh -o /srv/iso
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Benötigt Zugriff auf optische Laufwerke
DevicePolicy=closed
DeviceAllow=/dev/sr0 rw
DeviceAllow=/dev/cdrom rw

# Sicherheits-Einschränkungen
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/iso

[Install]
WantedBy=multi-user.target
