# disk2iso-lib/service/disk2iso.service
[Unit]
Description=disk2iso - Automatische ISO Erstellung von optischen Medien
After=multi-user.target
Wants=systemd-udevd.service
After=systemd-udevd.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/disk2iso -o /srv/iso
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Benötigt Zugriff auf optische Laufwerke
DevicePolicy=closed
DeviceAllow=/dev/sr0 rw
DeviceAllow=/dev/cdrom rw
DeviceAllow=block-* r

# Sicherheits-Einschränkungen (gelockert für Hardware-Zugriff)
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/srv/iso

[Install]
WantedBy=multi-user.target
